# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.error.log

[program:laravel-queue]
command=php /code/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
# numprocs=8
stdout_logfile=/var/log/laravel-queue.log
stderr_logfile=/var/log/laravel-queue.error.log
